# services/audio_stream.service
sudo tee /etc/systemd/system/audio_stream.service > /dev/null <<EOF
[Unit]
Description=Audio Streaming Server
After=network.target

[Service]
Type=simple
User=$REAL_USER
WorkingDirectory=$PROJECT_DIR
ExecStart=$VENV_PY $PROJECT_DIR/Audio/server.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
EOF
